# scripts/praier.service
[Unit]
Description=Praier - GitHub Pull Request Automation
After=network.target
Wants=network.target

[Service]
Type=simple
User=praier
Group=praier
WorkingDirectory=/opt/praier
ExecStart=/opt/praier/venv/bin/praier monitor
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=GITHUB_TOKEN=your_token_here
Environment=PRAIER_REPOSITORIES=owner/repo1,owner/repo2
Environment=PRAIER_POLL_INTERVAL=60
Environment=PRAIER_LOG_LEVEL=INFO

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/praier

[Install]
WantedBy=multi-user.target